# .config/systemd/user/albert.service
[Unit]
Description=albert

[Service]
ExecStart=/usr/bin/albert
Restart=always

[Install]
WantedBy=default.target
