# .config/systemd/user/nitrogen.service
[Unit]
Description=nitrogen

[Service]
ExecStart=/usr/bin/nitrogen --restore
Restart=on-failure

[Install]
WantedBy=default.target
