# .config/systemd/user/nitrogen.service
[Unit]
Description=nitrogen

[Service]
Type=oneshot
ExecStart=/usr/bin/nitrogen --restore

[Install]
WantedBy=graphical.target
