# dot_config/systemd/user/udiskie.service
[Unit]
Description=udiskie

[Service]
ExecStart=/usr/bin/udiskie
Restart=always

[Install]
WantedBy=default.target
