# systemd/.config/systemd/user/brightside.service
[Unit]
Description=Brightside

[Service]
ExecStart=/usr/bin/brightside
Restart=always

[Install]
WantedBy=default.target
