# home/dot_config/systemd/user/lxqt-policykit-agent.service
[Unit]
Description=lxqt-policykit-agent

[Service]
ExecStart=/usr/bin/lxqt-policykit-agent
Restart=always

[Install]
WantedBy=default.target
