# home/dot_config/systemd/user/rescuetime.service
[Unit]
Description=Rescuetime

[Service]
ExecStart=/usr/bin/rescuetime
Restart=always

[Install]
WantedBy=default.target
